# supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=Xvfb $DISPLAY -screen 0 1024x768x16